# tools/pulse-logsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.akgop.pulse-logsync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/newsfeed-pulse/tools/sync-logs.sh</string>
  </array>

  <!-- Every 15 minutes -->
  <key>StartInterval</key>
  <integer>900</integer>

  <!-- Also run on login -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/pulse-logsync.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/pulse-logsync.log</string>
</dict>
</plist>
